# chrony.conf
start on DEVICE_LOCK_DONE
stop on stopping xsession

# Restart automatically if it crashes
respawn
respawn limit 10 5

console output

script
    exec /usr/bin/chronyd -d
end script